# templates/default/systemd/alertmanager.service.erb
[Unit]
Description=Prometheus Alert Manager
After=network.target

[Service]
ExecStart=<%= node['prometheus']['alertmanager']['binary'] %> \
    -log.level=debug \
    -storage.path=<%= node['prometheus']['alertmanager']['storage.path'] %> \
    -config.file=<%= node['prometheus']['alertmanager']['config.file'] %> \
    -web.external-url=<%= node['prometheus']['alertmanager']['external-url'] %>
User=prometheus

[Install]
WantedBy=multi-user.target
